# ops/teamship-phase2-vm/newl-apps-auto-update.timer
[Unit]
Description=Check GitHub main for Newl Apps VM updates every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=newl-apps-auto-update.service

[Install]
WantedBy=timers.target
